# com.jyk.secretary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jyk.secretary</string>

    <key>WorkingDirectory</key>
    <string>/Users/jyk/Projects/secretary</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jyk/Projects/secretary/.venv/bin/python</string>
        <string>-m</string>
        <string>secretary.main</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/jyk/Projects/secretary/src:/Users/jyk/Projects/secretary</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jyk/Projects/secretary/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jyk/Projects/secretary/logs/launchd-stderr.log</string>
</dict>
</plist>
