# com.family.secretary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.family.secretary</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Update this path to your Python executable -->
        <string>/path/to/your/.venv/bin/python</string>
        <string>-m</string>
        <string>secretary.main</string>
    </array>
    <key>WorkingDirectory</key>
    <!-- Update this to your project path -->
    <string>/path/to/secretary</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/secretary/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/secretary/logs/launchd-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/path/to/secretary/src:/path/to/secretary/config/..</string>
    </dict>
</dict>
</plist>
